# bbs.conf
[program:pbbs1]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/pbbs1
autostart=true
autorestart=true